# systemd/systemd-ask-password-gui.path
[Unit]
Description=UI for boot time passwords (path)
After=wait_for_keymaster.service droid-late-start.service late-mount-pre.target
Requires=wait_for_keymaster.service droid-late-start.service late-mount-pre.target
Before=late-mount.target
DefaultDependencies=no
Conflicts=shutdown.target

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
