# systemd/systemd-ask-password-gui.path
[Unit]
Description=UI for boot time passwords (path)
After=late-mount-pre.target
Requires=late-mount-pre.target
Before=late-mount.target
DefaultDependencies=no
Conflicts=shutdown.target
ConditionPathExists=!/run/systemd/boot-status/ACT_DEAD
ConditionPathExists=!/run/systemd/boot-status/TEST

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
